# templates/default/init/minecraft.conf.erb
# Minecraft Server
#
# Note: This file is managed by Chef. Any
# manual changes WILL be reverted.

description "Minecraft Server"

start on runlevel [2345]
stop on runlevel [!2345]

exec su -c '<%= node[:minecraft][:base_dir] %>/server/server.sh' minecraft
